# guardian/cron.d/01-rotor
#
# Rotate/remove backups
55 1 * * * /home/guardian/bin/rotor.sh
