# guardian/cron.d/05-backup
#
# Daily backup of SQLite databases in /var/data
0 2 * * * /home/guardian/bin/backup.sh
